# project/conf/supervisor.conf
[supervisord]
nodaemon=true

[program:mongo]
command=mongod --logpath=/dev/null --bind_ip_all -f /etc/mongod.conf

[program:cron]
command=cron -f -L 2
